# deployment/crontab/shxreader
SHELL=/bin/bash
SHX_PATH=/home/pi/SHXReader
0 12 * * * pi python ${SHX_PATH}/shxreader.py
