# deployment/crontab/shxreader
SHELL=/bin/bash
SHX_PATH=/home/pi/SHXReader
0 12 * * * pi cd ${SHX_PATH} && python ${SHX_PATH}/shxreader.py > /home/pi/shxreader.log 2>&1
0 12 * * * pi cd ${SHX_PATH} && python ${SHX_PATH}/shxcrawler.py > /home/pi/shxcrawler.log 2>&1
